# backend/vpn-api.service
[Unit]
Description=VPN Account Management API
Documentation=https://github.com/your-repo/backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/scAUTO/backend
ExecStart=/root/scAUTO/backend/vpn-api
Restart=always
RestartSec=3
Environment=PORT=3000
Environment=be51e916-7380-47c4-98dd-d05c3d26ecd2

[Install]
WantedBy=multi-user.target